# ai-fconfig/etc-systemd-system-fbpanel_update.path
[Unit]
Description=Check update logo fbpanel

[Path]
PathChanged=/usr/share/fbpanel/images/logo.png

Unit=fbpanel_update.service

[Install]
WantedBy=multi-user.target
