# ai-fconfig/efi-services/refind_update.path
[Unit]
Description=Check update rEFInd

[Path]
PathChanged=/usr/share/refind/refind_x64.efi

Unit=refind_update.service

[Install]
WantedBy=multi-user.target
